# deployments/runtimes/fluent/supervisor/conf.d/laravel-queue-worker.ini
[program:laravel-queue-worker]
command=php /app/artisan queue:work
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/queue-worker.err.log
stdout_logfile=/var/log/supervisor/queue-worker.out.log